# files/system/etc/profile.d/99-ssh.sh
# Set SSH_AUTH_SOCK to use the socket in XDG_RUNTIME_DIR
export SSH_AUTH_SOCK="${XDG_RUNTIME_DIR}/ssh-agent.socket"